# debian/postinst
#!/bin/sh

set -e

. /etc/default/openmediavault
. /usr/share/openmediavault/scripts/helper-functions

OMV_COUCHPOTATO_APP_PATH="/opt/couchpotato"
OMV_COUCHPOTATO_DATA_PATH="/var/opt/couchpotato"
OMV_COUCHPOTATO_USER="couchpotato"

migrate_data_location()
{
    old_data_path="/home/couchpotato/.couchpotato"

    if [ -d /home/couchpotato/.couchpotato ]; then
        # Migrate old data to new location.
        su -s /bin/sh -c "cp -ra /home/couchpotato/.couchpotato/* ${OMV_COUCHPOTATO_DATA_PATH}" $OMV_COUCHPOTATO_USER

        # Remove couchpotato home directory.
        if [ -d /home/couchpotato ]; then
            rm -rf /home/couchpotato
        fi
    fi
}

case "$1" in
    configure)
        # Set the default settings of the service package to those expected by
        # OpenMediaVault.
        if [ -z "${2}" ]; then
            systemctl stop couchpotato.service >/dev/null || true
            systemctl disable couchpotato.service >/dev/null || true
        fi

        SERVICE_XPATH_NAME="couchpotato"
        SERVICE_XPATH="/config/services/${SERVICE_XPATH_NAME}"



        if omv_config_exists "${SERVICE_XPATH}"; then
            ENABLE="$(omv_config_get "${SERVICE_XPATH}/enable")"
            SHOWTAB="$(omv_config_get "${SERVICE_XPATH}/show_tab")"
            BRANCH="$(omv_config_get "${SERVICE_XPATH}/branch")"
            REPO="$(omv_config_get "${SERVICE_XPATH}/repo")"
        else
            ENABLE="0"
            SHOWTAB="0"
            BRANCH="Please choose a branch"
            REPO="Please choose a repository"
        fi

        if omv_config_exists "${SERVICE_XPATH}/username"; then
            USERNAME="$(omv_config_get "${SERVICE_XPATH}/username")"
            USERGRP="$(omv_config_get "${SERVICE_XPATH}/usersgrp")"
            UMASK="$(omv_config_get "${SERVICE_XPATH}/umask")"
        else
            USERNAME="couchpotato"
            USERGRP="0"
            UMASK="000"
        fi

        if omv_config_exists "${SERVICE_XPATH}/repo2"; then
            BRANCH2="$(omv_config_get "${SERVICE_XPATH}/branch2")"
            REPO2="$(omv_config_get "${SERVICE_XPATH}/repo2")"
            SSL="$(omv_config_get "${SERVICE_XPATH}/ssl")"
            PPASS="$(omv_config_get "${SERVICE_XPATH}/ppass")"
            INSTAL="$(omv_config_get "${SERVICE_XPATH}/newinstance")"
            RUN="$(omv_config_get "${SERVICE_XPATH}/newinstenable")"
            PORT="$(omv_config_get "${SERVICE_XPATH}/port")"
        else
            SSL="0"
            PPASS="0"
            BRANCH2="Please choose a branch"
            REPO2="Please choose a repository"
            INSTAL="0"
            RUN="0"
            PORT="5050"
        fi

        omv_config_delete "${SERVICE_XPATH}"

        omv_config_add_element "/config/services" "${SERVICE_XPATH_NAME}"
        omv_config_add_element "${SERVICE_XPATH}" "enable" $ENABLE
        omv_config_add_element "${SERVICE_XPATH}" "repo" $REPO
        omv_config_add_element "${SERVICE_XPATH}" "branch" $BRANCH
        omv_config_add_element "${SERVICE_XPATH}" "repo2" $REPO2
        omv_config_add_element "${SERVICE_XPATH}" "branch2" $BRANCH2
        omv_config_add_element "${SERVICE_XPATH}" "show_tab" $SHOWTAB
        omv_config_add_element "${SERVICE_XPATH}" "username" $USERNAME
        omv_config_add_element "${SERVICE_XPATH}" "usersgrp" $USERGRP
        omv_config_add_element "${SERVICE_XPATH}" "umask" $UMASK
        omv_config_add_element "${SERVICE_XPATH}" "port" $PORT
        omv_config_add_element "${SERVICE_XPATH}" "newinstance" $INSTAL
        omv_config_add_element "${SERVICE_XPATH}" "newinstenable" $RUN
        omv_config_add_element "${SERVICE_XPATH}" "ssl" $SSL
        omv_config_add_element "${SERVICE_XPATH}" "PPASS" $PPASS

        if ! getent passwd couchpotato >/dev/null 2>&1; then
            adduser --quiet \
                    --system \
                    --group \
                    --no-create-home \
                    --disabled-password \
                    $OMV_COUCHPOTATO_USER
        fi

        chown $OMV_COUCHPOTATO_USER:$OMV_COUCHPOTATO_USER $OMV_COUCHPOTATO_APP_PATH
        chown $OMV_COUCHPOTATO_USER:$OMV_COUCHPOTATO_USER $OMV_COUCHPOTATO_DATA_PATH

        if dpkg --compare-versions "${2}" lt-nl "1.1"; then
            migrate_data_location
        fi

        # Activate package triggers. These triggers are only set during the
        # package installation.
        dpkg-trigger update-fixperms
        dpkg-trigger update-locale
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
